# smta.conf
[program:smta]
directory=/root/smta2
command=python smta.py
autostart=true
autorestart=true
